# apache2-svn/authz/build/site/etc/svnauth.cron
*/5 * * * *     www-data        /usr/bin/python /var/lib/site/script/svnauth.py /var/lib/site/etc/svnauth.json
